# recipes-multimedia/aml-tvserver/files/tvserver.service
[Unit]
Description=TV Server Startup
After=binder.service

[Service]
User=root
ExecStart=/usr/bin/tvservice

[Install]
WantedBy=multi-user.target
