# recipes-core/aml-hdcp/files/aml_hdcp.service
[Unit]
Description=hdcp Startup
After=tee-supplicant.service

[Service]
User=root
ExecStart=/usr/bin/tee_hdcp

[Install]
WantedBy=multi-user.target
