# recipes-core/aml-hdcp/files/aml_hdcp.service
[Unit]
Description=hdcp Startup
After=aml_key_inject.service

[Service]
User=root
ExecStart=/usr/bin/tee_hdcp

[Install]
WantedBy=multi-user.target
