# recipes-rtos/aml-rtos/files/load_dsp.service
[Unit]
Description=DSP firmware load

[Service]
Type=forking
User=root
ExecStart=/usr/bin/load_dsp.sh start

[Install]
WantedBy=multi-user.target
